# firstboot
#!/sbin/openrc-run
command=/sbin/openrc-firstboot
command_args=""

name="openrc-firstboot"
description="Firstboot helper for new users."
